# ark@.timer
[Unit]
Description=Boot delay for ARK server

[Timer]
OnBootSec=1min

[Install]
WantedBy=timers.target